# etc/flux-prolog@.service.in
[Unit]
Description=Prolog for Flux job %I
CollectMode=inactive-or-failed

[Service]
Type=oneshot
EnvironmentFile=-@X_RUNSTATEDIR@/flux-prolog@%I.env
ExecStart=@X_LIBEXECDIR@/flux/flux-run-system-scripts prolog
ExecStopPost=-rm -f @X_RUNSTATEDIR@/flux-prolog@%I.env
