# etc/flux.service.in
[Unit]
Description=Flux message broker
Wants=munge.service

[Service]
TimeoutStopSec=90
KillMode=mixed
ExecStart=/bin/bash -c '\
  XDG_RUNTIME_DIR=/run/user/$UID \
  DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/$UID/bus \
  @X_BINDIR@/flux broker \
  --config-path=@X_SYSCONFDIR@/flux/system/conf.d \
  -Scron.directory=@X_SYSCONFDIR@/flux/system/cron.d \
  -Stbon.fanout=256 \
  -Srundir=@X_RUNSTATEDIR@/flux \
  -Slocal-uri=local://@X_RUNSTATEDIR@/flux/local \
  -Slog-stderr-level=6 \
  -Slog-stderr-mode=local \
  -Scontent.backing-path=@X_LOCALSTATEDIR@/lib/flux/content.sqlite \
  -Sbroker.rc2_none \
  -Sbroker.quorum=0 \
  -Sbroker.quorum-timeout=none \
'
SyslogIdentifier=flux
ExecReload=@X_BINDIR@/flux config reload
Restart=always
RestartSec=5s
User=flux
Group=flux
RuntimeDirectory=flux
RuntimeDirectoryMode=0755
PermissionsStartOnly=true
ExecStartPre=-/bin/mkdir -p @X_LOCALSTATEDIR@/lib/flux
ExecStartPre=-/bin/chown flux:flux @X_LOCALSTATEDIR@/lib/flux
ExecStartPre=-/bin/chmod 0700 @X_LOCALSTATEDIR@/lib/flux
ExecStartPre=/usr/bin/loginctl enable-linger flux
ExecStartPre=bash -c 'systemctl start user@$(id -u flux).service'

#
# Delegate cgroup control to user flux, so that systemd doesn't reset
#  cgroups for flux initiated processes, and to allow (some) cgroup
#  manipulation as user flux.
#
Delegate=yes

[Install]
WantedBy=multi-user.target
