# etc/flux.service.in
[Unit]
Description=Flux message broker

[Service]
Environment=FLUX_USERDB_OPTIONS=--default-rolemask=user
ExecStart=@X_BINDIR@/flux start -o,-Sbroker.rundir=@X_RUNSTATEDIR@/flux,-Ssession-id=%H sleep inf
User=flux
Group=flux
RuntimeDirectory=flux
RuntimeDirectoryMode=0755
#
# Delegate cgroup control to user flux, so that systemd doesn't reset
#  cgroups for flux initiated processes, and to allow (some) cgroup
#  manipulation as user flux.
#
Delegate=yes

[Install]
WantedBy=multi-user.target
