# etc/flux.service.in
[Unit]
Description=Flux message broker

[Service]
Environment=FLUX_USERDB_OPTIONS=--default-rolemask=user
TimeoutStopSec=90
KillMode=mixed
ExecStart=@X_BINDIR@/flux broker \
  -Srundir=@X_RUNSTATEDIR@/flux \
  -Slocal-uri=local://@X_RUNSTATEDIR@/flux/local \
  -Sboot.method=config \
  -Slog-stderr-level=6 \
  -Scontent.backing-path=@X_LOCALSTATEDIR@/lib/flux/content.sqlite \
  -Sbroker.rc2_none
User=flux
Group=flux
RuntimeDirectory=flux
RuntimeDirectoryMode=0755
#
# Delegate cgroup control to user flux, so that systemd doesn't reset
#  cgroups for flux initiated processes, and to allow (some) cgroup
#  manipulation as user flux.
#
Delegate=yes

[Install]
WantedBy=multi-user.target
