# etc/flux-housekeeping@.service.in
[Unit]
Description=Housekeeping for Flux job %I
CollectMode=inactive-or-failed

[Service]
Type=oneshot
EnvironmentFile=-@X_RUNSTATEDIR@/flux-housekeeping@%I.env
ExecStart=@X_LIBEXECDIR@/flux/flux-run-system-scripts housekeeping
ExecStopPost=-rm -f @X_RUNSTATEDIR@/flux-housekeeping@%I.env
ExecStopPost=-sh -c '\
    if test "$SERVICE_RESULT" != "success"; then \
        message="housekeeping@%I ${SERVICE_RESULT:-failure}"; \
        if test "${EXIT_CODE}${EXIT_STATUS}"; then \
            message="$message: $EXIT_CODE $EXIT_STATUS"; \
        fi; \
        flux resource drain $(flux getattr rank) $message; \
    fi \
'
