# etc/flux.service.in
[Unit]
Description=Flux message broker

[Service]
Environment=FLUX_USERDB_OPTIONS=--default-rolemask=user
ExecStart=@X_BINDIR@/flux broker -Srundir=@X_RUNSTATEDIR@/flux -Sboot.method=config -Sboot.config_file=@X_SYSCONFDIR@/flux/conf.d/boot.conf sleep inf
User=flux
Group=flux
RuntimeDirectory=flux
RuntimeDirectoryMode=0755
#
# Delegate cgroup control to user flux, so that systemd doesn't reset
#  cgroups for flux initiated processes, and to allow (some) cgroup
#  manipulation as user flux.
#
Delegate=yes

[Install]
WantedBy=multi-user.target
